# infra/systemd/daybook.service
# /etc/systemd/system/daybook.service
# Runs the daybook container stack via docker compose, owned by daybookuser.
# Mirrors the otuburu service model. Enable with:
#   sudo systemctl enable --now daybook
[Unit]
Description=Daybook (multi-tenant daily reporting)
Documentation=https://github.com/filatei/daybook
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
User=daybookuser
Group=daybookuser
WorkingDirectory=/opt/daybook/backend
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose up -d --force-recreate --remove-orphans
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
